# from_scratch_checkpoint_on_jazz_server/app/texting-coding-llm.service
[Unit]
Description=Texting Coding Local LLM API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/texting-coding-model/app
ExecStart=/opt/texting-coding-model/venv/bin/uvicorn app:app --host 127.0.0.1 --port 18080
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
